# config/com.paireyewear.labassistant.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>Label</key>
    <string>com.paireyewear.labassistant.gateway</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/npx</string>
        <string>tsx</string>
        <string>index.ts</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/Shared/lab_assistant/server/data/gateway-error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/Shared/lab_assistant/server/data/gateway.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/Shared/lab_assistant/gateway</string>
</dict>
</plist>
